# templates/etc/supervisor/conf.d/dashboard.conf
[program:dashboard]
command = /usr/local/bin/gunicorn --config /etc/archivematica/dashboard.gunicorn-config.py wsgi:application

