# templates/etc/supervisor/conf.d/mcp-server.conf
[program:MCPServer]
environment = CONF=/etc/archivematica/MCPServer
command = /usr/lib/archivematica/MCPServer/archivematicaMCP.py
user = archivematica
