# templates/etc/supervisor/conf.d/mcp-client.conf
[program:MCPClient]
environment = CONF=/etc/archivematica/MCPClient,HOME=/var/lib/archivematica/
command = /usr/lib/archivematica/MCPClient/archivematicaClient.py
user = archivematica
