# files/supervisor/fits.conf
[program:FITS]
user = archivematica
command = fits-ngserver.sh /usr/share/maven-repo/com/martiansoftware/nailgun-server/0.9.1/nailgun-server-0.9.1.jar
