# files/supervisor/storage-service.conf
[program:storage-service]
command = /usr/share/python/archivematica-storage-service/bin/gunicorn --config /etc/archivematica/storage-service.gunicorn-config.py storage_service.wsgi:application
